# server/systemd/ppocrv5-http.service
[Unit]
Description=PP-OCRv5 HTTP Server (JK-LLM-OCR)
After=network.target
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Type=simple

# IMPORTANT:
# - Keep the server bound to 127.0.0.1 and access it via SSH tunnel.
# - Update paths below to match your actual install directory.
WorkingDirectory=/opt/jk-llm-ocr/server

Environment=HOST=127.0.0.1
Environment=PORT=50000
Environment=DEVICE=cpu
Environment=PADDLE_PDX_DISABLE_MODEL_SOURCE_CHECK=True

ExecStart=/opt/jk-llm-ocr/server/scripts/start_ppocrv5_server.sh

Restart=on-failure
RestartForceExitStatus=SIGTERM
RestartSec=3

[Install]
WantedBy=multi-user.target
